# bot.conf
[program: mortgage]
directory=/home/mortgage
command=python app.py

autostart=true
autorestart=true
stderr_logfile=/home/mortgage/logfile_err.log
stdout_logfile=/home/mortgage/logfile.log
