# bot.conf
[program: 420]
directory=/home/420
command=python app.py

autostart=true
autorestart=true
stderr_logfile=/home/420/logfile_err.log
stdout_logfile=/home/420/logfile.log
