# money.conf
[program: money]
directory=/home/money
command=python app.py

autostart=true
autorestart=true
stderr_logfile=/home/money/logfile_err.log
stdout_logfile=/home/money/logfile.log
